# scripts/com.polymarket.weatherbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.polymarket.weatherbot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/plux96/polymarket-weather-bot/venv/bin/python</string>
        <string>/Users/plux96/polymarket-weather-bot/bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/plux96/polymarket-weather-bot</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/plux96/polymarket-weather-bot/bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/plux96/polymarket-weather-bot/bot_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
